# meta-boot2qt/recipes-qt/automotive/qtivi/ivi-services.target
[Unit]
Description=IVI services
After=systemd-user-sessions.service
